# squid/supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=/opt/run_squid.sh
redirect_stderr = true

[program:tt-server]
command=/opt/tt-server.sh
redirect_stderr = true
